# cmd/scaletaild/scaletail-online.target
[Unit]
Description=ScaleTail is online
Requires=scaletail-wait-online.service
After=scaletail-wait-online.service
